# GoFBot/build/stageStatic/home/gofbot/GoFBot/crontab
45 17 * * 1,3,5 /usr/bin/sudo /usr/sbin/service gofbot start

